# deploy/vps/mapmapmaps-autopull.timer
[Unit]
Description=Poll GitHub for MapMapMaps updates

[Timer]
OnBootSec=2min
OnUnitActiveSec=2min
AccuracySec=30s

[Install]
WantedBy=timers.target
